# init/rigctl-http.service
[Unit]
Description=HTTP interface service for rigctld

[Service]
Type=simple
User=root
Group=root
ExecStart=rigctl-http
